# whatsapp-bot.service
[Unit]
Description=WhatsApp Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/your/app
ExecStart=/usr/bin/node app.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PWD=/path/to/your/app

# Logs
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=whatsapp-bot

# Seguridad
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

# INSTRUCCIONES DE USO:
# 1. Copia este archivo a /etc/systemd/system/whatsapp-bot.service
# 2. Cambia /path/to/your/app por la ruta real
# 3. Cambia 'ubuntu' por tu usuario
# 4. Ejecuta: sudo systemctl daemon-reload
# 5. Ejecuta: sudo systemctl enable whatsapp-bot
# 6. Ejecuta: sudo systemctl start whatsapp-bot
